# packer/vault/ansible/roles/vault_node/templates/vault.service.j2
[Unit]
Description=HashiCorp Vault
Documentation=https://www.vaultproject.io/docs
Requires=network-online.target vault-firstboot.service
After=network-online.target vault-firstboot.service
ConditionFileNotEmpty={{ vault_node_config_dir }}/vault.hcl

[Service]
User={{ vault_node_user }}
Group={{ vault_node_group }}
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes
SecureBits=keep-caps
AmbientCapabilities=CAP_IPC_LOCK
CapabilityBoundingSet=CAP_SYSLOG CAP_IPC_LOCK
NoNewPrivileges=yes
# -config points at the DIRECTORY -- vault server merges all *.hcl files
# inside (vault.hcl + optional seal-transit.hcl). The seal-transit.hcl
# drop-in is delivered post-clone by Terraform on cluster members
# (vault-1/2/3) for Phase 0.D.5.5 transit auto-unseal. vault-transit
# itself never receives one (stays in shamir mode; it IS the unseal
# key custodian for the others).
ExecStart=/usr/local/bin/vault server -config={{ vault_node_config_dir }}/
ExecReload=/bin/kill --signal HUP $MAINPID
KillMode=process
KillSignal=SIGINT
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
StartLimitInterval=60
StartLimitBurst=3
LimitNOFILE=65536
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
